# deploy/dev-fleet.service
[Unit]
Description=Mission Control fleet — autonomous build agents (supervisor loop)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fleet
WorkingDirectory=/home/fleet/agent-dashboard
ExecStart=/home/fleet/agent-dashboard/loop.sh
# claude/gh/git/npm must be on PATH for the worker; HOME is set automatically from User=.
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/home/fleet/.local/bin
# Config comes from config.env + config.local.env (sourced by lib.sh). Start conservatively:
Environment=MAX_WORKERS=1
Environment=XDG_RUNTIME_DIR=/run/user/1000
Restart=always
RestartSec=5

# ── hardening ──
NoNewPrivileges=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictSUIDSGID=true
# NOTE: this isolates the ORCHESTRATOR, not the agent step. The 'claude -p' build agent runs with
# --dangerously-skip-permissions and must run in a separate container (deploy/run-agent-sandbox.sh)
# without $HOME secrets + with an egress allowlist before you run an untrusted backlog 24/7. See deploy/FASE-0-CHECKLIST.md.

[Install]
WantedBy=multi-user.target
