# deploy/mission-control-dashboard.service
[Unit]
Description=Mission Control dashboard (Next.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fleet
WorkingDirectory=/home/fleet/agent-dashboard/mission-control
# Bind to localhost only — reachable via Tailscale Serve, never the public interface.
ExecStart=/usr/bin/npm run start -- -H 127.0.0.1 -p 3000
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Restart=always
RestartSec=5

# ── hardening ──
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
