# deploy/mission-control-watchdog.timer
[Unit]
Description=Run the Mission Control fleet watchdog every 2 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=30s
Unit=mission-control-watchdog.service

[Install]
WantedBy=timers.target
