# deploy/morning-digest.timer
[Unit]
Description=Run the Mission Control morning digest daily at 07:15

[Timer]
OnCalendar=*-*-* 07:15:00
Persistent=true
AccuracySec=1min
Unit=morning-digest.service

[Install]
WantedBy=timers.target
